# android9/overlay/usr/share/upstart/sessions/clear-gst-registry.conf
description "Clear GStreamer registry cache"
author "Alfred Neumayer <dev.beidl@gmail.com>"

start on starting mediascanner-2.0

script
    rm -rf $HOME/.cache/gstreamer-1.0
end script
